# misc/systemd/gutdeb-autopublish@.service
[Unit]
Description=Gutdeb autopublish service
After=network.target

[Service]
User=%i
ExecStart=/usr/bin/gutdeb autopublish /app/debs/*
